# system_files/shared/usr/lib/systemd/system/flatpak-appstream-refresh.service
[Unit]
Description=Refresh Flatpak remote metadata
Documentation=https://docs.projectbluefin.io
After=network-online.target flatpak-system-helper.service
Wants=network-online.target
ConditionPathExists=/usr/bin/flatpak
ExecCondition=bash -c 'v=$(busctl get-property org.freedesktop.NetworkManager /org/freedesktop/NetworkManager org.freedesktop.NetworkManager Metered 2>/dev/null | cut -d" " -f2); [ "$v" != "1" ] && [ "$v" != "2" ]'
StartLimitIntervalSec=600
StartLimitBurst=3

[Service]
Type=oneshot
ExecStart=/usr/bin/flatpak update --appstream --system --noninteractive
RemainAfterExit=yes
Restart=on-failure
RestartSec=60

[Install]
WantedBy=graphical.target
