# system_files/shared/etc/profile.d/uwelcome.sh
#!/usr/bin/env bash

rm -rf ~/.config/no-show-user-motd || true # Made to be left here only for the time every Bluefin switch to uwelcome
uwelcome
